# examples/weather-cleanup.timer
[Unit]
Description=Remove outdated weather forecasts

[Timer]
OnCalendar=daily
AccuracySec=1h

[Install]
WantedBy=timers.target
